# 15. Базы данных. Продолжение/hw-15/hw15-3.service
[Unit]
Description="Сервис для домашки 15-3 - бэкапы mysql"

[Service]
ExecStart=/usr/local/bin/backup.sh

[Install]
WantedBy=multi-user.target